# crontab
*/15 0-7,13-23 * * * root cd /usr/src/ebay && /usr/local/bin/python run.py ebay --configfile=config.json --loglevel=INFO
